# biga.service.example
[Unit]
Description=Big A Angels Tracker
After=network-online.target multi-user.target
Wants=network-online.target
# Only start when user WiFi credentials have been provisioned via the portal.
ConditionPathExists=/etc/biga/wifi_creds.json

[Service]
Type=simple
User=root
WorkingDirectory=/home/pi/BigA
ExecStart=/usr/local/bin/biga-start.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target