# biga.service.example
[Unit]
Description=Big A Angels Tracker
After=network-online.target multi-user.target biga-connectivity.service
Wants=network-online.target
Requires=biga-connectivity.service
# Scoreboard runs when WiFi is provisioned and not in add-network / offline mode.
ConditionPathExists=/etc/biga/wifi_creds.json
ConditionPathExists=!/etc/biga/provisioning_active

[Service]
Type=simple
User=root
WorkingDirectory=/home/pi/BigA
# Quick OTA attempt before launch (non-fatal). biga-update-boot.timer retries
# ~3 min after boot when WiFi/DNS are usually stable.
ExecStartPre=-/bin/bash -c 'BIGA_UPDATE_CONTEXT=boot-pre /home/pi/BigA/scripts/update_biga.sh'
ExecStart=/usr/local/bin/biga-start.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target