# scripts/biga-connectivity.service
[Unit]
Description=BigA WiFi boot sync (NM profiles from saved networks)
After=NetworkManager.service network-online.target
Wants=network-online.target
Before=biga.service biga-portal.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/python3 /home/pi/BigA/scripts/check_wifi_connectivity.py
StandardOutput=append:/var/log/biga-connectivity.log
StandardError=append:/var/log/biga-connectivity.log

[Install]
WantedBy=multi-user.target
