# portal/biga-portal.service
[Unit]
Description=BigA WiFi Provisioning Portal
After=network.target NetworkManager.service biga-firstboot.service
# Only runs when no user WiFi credentials exist (AP / factory-reset state).
ConditionPathExists=!/etc/biga/wifi_creds.json

[Service]
Type=simple
User=root
WorkingDirectory=/home/pi/BigA/portal
ExecStartPre=/usr/bin/nmcli con up biga-ap
ExecStart=/usr/bin/python3 /home/pi/BigA/portal/portal.py
ExecStopPost=/usr/bin/nmcli con down biga-ap
Restart=on-failure
RestartSec=5
StandardOutput=append:/var/log/biga-portal.log
StandardError=append:/var/log/biga-portal.log

[Install]
WantedBy=multi-user.target
